# deploy/linux/systemd/meeting-digest-bot-knowledge-import.timer.example
[Unit]
Description=Run MeetingDigestBot Notion back-import checks twice daily

[Timer]
OnCalendar=*-*-* 10,18:00:00 Europe/Kyiv
Persistent=true
Unit=meeting-digest-bot-knowledge-import.service

[Install]
WantedBy=timers.target
