# deploy/linux/systemd/meeting-digest-bot-poller.service.example
[Unit]
Description=MeetingDigestBot Telegram poller
Wants=network-online.target
After=network-online.target meeting-digest-bot-api.service

[Service]
Type=simple
User=YOUR_LINUX_USER
Group=YOUR_LINUX_USER
WorkingDirectory=/opt/meeting-digest-bot
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/meeting-digest-bot/.env
ExecStart=/opt/meeting-digest-bot/.venv/bin/python -m meeting_digest_bot poll-telegram --drop-pending
Restart=always
RestartSec=5
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
