# deploy/linux/systemd/meeting-digest-bot-knowledge-import.timer.example
[Unit]
Description=Run MeetingDigestBot Notion back-import checks periodically

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
Persistent=true
Unit=meeting-digest-bot-knowledge-import.service

[Install]
WantedBy=timers.target
