# deploy/linux/systemd/task-extractor-bot-poller.service.example
[Unit]
Description=Task Extractor Telegram poller
Wants=network-online.target
After=network-online.target meeting-digest-bot-api.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/meeting-digest-bot
Environment=PYTHONUNBUFFERED=1
UMask=0002
EnvironmentFile=/opt/meeting-digest-bot/.env
ExecStart=/opt/meeting-digest-bot/.venv/bin/python -m meeting_digest_bot poll-task-extractor --drop-pending
Restart=always
RestartSec=5
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
